# dockerfiles/libvirt/overlay/etc/systemd/system/teleport.service.d/override.conf
[Service]
EnvironmentFile=-/etc/default/teleport

ExecStart=
ExecStart=/usr/bin/teleport start --roles=node --nodename=db --config $CONFIG
SyslogIdentifier=teleport-ssh-client
