# dockerfiles/libvirt/overlay/etc/systemd/system/teleport.service.d/override.conf
[Service]
EnvironmentFile=-/etc/default/teleport

ExecStart=
ExecStart=/usr/bin/teleport start --roles=node --token=internalxyz --auth-server=teleport:3025
SyslogIdentifier=teleport-ssh-client
